# config/systemd/user/backup-local-onera.timer
[Unit]
Description=Run `backup-local-onera.service` daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
